# supervisord/plex-api.conf
[program:plex-api]
autorestart = true
autostart = true
startsecs = 0
directory = /root/plex-api/
command = /usr/bin/python2.7 app.py
startretries = 1
environment=PYTHONPATH="%(here)s"
